# nexuspad/.config/systemd/user/battery-alert.service
[Unit]
Description=Battery status notification

[Service]
Type=oneshot
ExecStart=%h/.local/bin/battery-alert
