# nexuspad/.config/systemd/user/battery-alert.timer
[Unit]
Description=Check battery status every 1 minute

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
Persistent=true

[Install]
WantedBy=timers.target